# init.d/openrc/jacklistener.in
#!/sbin/runscript

PIDFILE=/var/run/jacklistenerd.pid
DAEMON="@sbindir/jacklistener-runscript"

depend() {
	need localmount
	after bootmisc modules isapnp coldplug hotplug
}

start() {
	ebegin "Starting Jack Listener Daemon"
	start-stop-daemon --start --quiet --pidfile $PIDFILE --exec "$DAEMON" -- --quiet --pidfile $PIDFILE --daemon
	eend $?
}

stop() {
	ebegin "Stopping Jack Listener Daemon"
	start-stop-daemon --stop --quiet --pidfile $PIDFILE
	rm -f $PIDFILE
	eend $?
}
