# init.d/openrc/jacklistener
#!/sbin/runscript

PIDFILE=/var/run/jacklistenerd.pid

depend() {
	need localmount
	after bootmisc modules isapnp coldplug hotplug
}

start() {
	ebegin "Starting Jack Listener Daemon"
	DEVICES=
	cd /sys/class/input
	for evdev in event*; do
		PHYS="`cat "$evdev/device/phys"`"
		if [ "$PHYS" = ALSA ]; then
			DEVICES="$DEVICES /dev/input/$evdev"
		fi
	done
	cd "$OLDPWD"
	start-stop-daemon --start --quiet --pidfile $PIDFILE --exec /usr/sbin/jacklistenerd -- --quiet --pidfile $PIDFILE --daemon $DEVICES
	eend $?
}

stop() {
	ebegin "Stopping Jack Listener Daemon"
	start-stop-daemon --stop --quiet --pidfile $PIDFILE
	rm -f $PIDFILE
	eend $?
}
