# init.d/jacklistener.service.in
[Unit]
Description=Jack Listener Daemon
After=sysinit.target

[Service]
Type=simple
ExecStart=@sbindir@/jacklistener-runscript

[Install]
WantedBy=multi-user.target
