# init.d/jacklistener-openrc
#!/sbin/runscript

depend() {
	need localmount
	after bootmisc modules isapnp coldplug hotplug
}

start() {
	ebegin "Starting Jack Listener Daemon"
	start-stop-daemon --start --quiet --background --pidfile /var/run/jacklistenerd.pid --make-pidfile --exec /usr/sbin/jacklistenerd
	eend $?
}

stop() {
	ebegin "Stopping Jack Listener Daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/jacklistenerd.pid
	rm -f /var/run/jacklistenerd.pid
	eend $?
}
